# fio.service
[Unit]
Description=FIO server

[Service]
ExecStart=/usr/bin/fio --server

[Install]
WantedBy=multi-user.target

